# root/etc/services.d/budge/run
#!/usr/bin/with-contenv bash

cd /app/budge/backend

exec \
    s6-setuidgid abc /usr/bin/npm run start
